# eye_detector_v2/rootfs/etc/services.d/eye_detector/run
#!/usr/bin/with-contenv bash
exec python3 /app/eye_detector.py
